# journal/extensions/systemd/journal.service
[Unit]
Description=Journal background watcher
After=graphical-session.target

[Service]
ExecStart=/home/fira/Documents/f/journal/run.sh
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
